# video-monitor.service
# systemd 服务文件模板
# 将此文件复制到 /etc/systemd/system/video-monitor.service
# 然后根据实际情况修改路径和用户名

[Unit]
Description=Video Monitor Service - 视频监控服务
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=yusteven
Group=yusteven
WorkingDirectory=/home/yusteven/boluo/pyt
ExecStart=/bin/bash /home/yusteven/boluo/pyt/monitor.sh start-foreground
ExecStop=/bin/bash /home/yusteven/boluo/pyt/monitor.sh stop
ExecReload=/bin/bash /home/yusteven/boluo/pyt/monitor.sh restart

# 重启设置
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# 环境变量
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/usr/local/bin:/usr/bin:/bin

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/yusteven/boluo/pyt

# 日志设置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=video-monitor

[Install]
WantedBy=multi-user.target

# 安装说明:
# 1. 复制此文件到 /etc/systemd/system/video-monitor.service
# 2. 修改 User, Group, WorkingDirectory 路径（如需要）
# 3. 重新加载 systemd: sudo systemctl daemon-reload
# 4. 启用服务: sudo systemctl enable video-monitor
# 5. 启动服务: sudo systemctl start video-monitor
# 6. 查看状态: sudo systemctl status video-monitor
# 7. 查看日志: sudo journalctl -u video-monitor -f
#
# 注意: 服务使用 monitor.sh start-foreground 启动，确保脚本有执行权限
