# video-monitor.service
[Unit]
Description=Video Monitor Service (Direct Python)
After=network-online.target
Wants=network-online.target

[Service]
# 指定运行服务的用户和组
User=yusteven
Group=yusteven

# 设置工作目录
WorkingDirectory=/home/yusteven/boluo/pyt

# 直接指向虚拟环境中的 Python 解释器
ExecStart=/home/yusteven/.pyenv/versions/3.14.2/bin/python /home/yusteven/boluo/pyt/monitor.py

# 如果 monitor.py 需要从 .env 文件读取变量，systemd 可以直接加载它
EnvironmentFile=-/home/yusteven/boluo/pyt/data/.env

# 确保 Python 输出直接进入 journal 日志，而不被缓存
Environment=PYTHONUNBUFFERED=1
# 将虚拟环境的 bin 目录加入 PATH，确保脚本内调用 yt-dlp 等命令时能找到
Environment=PATH=/home/yusteven/boluo/pyt/.venv/bin:/usr/local/bin:/usr/bin:/bin

# 重启策略
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

# --- 安全与沙盒设置 ---
# 禁止获取新特权
NoNewPrivileges=true
# 隔离临时目录
PrivateTmp=true
# 设置文件系统为只读，仅允许写入指定的路径
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/yusteven/boluo/pyt

# 日志配置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=video-monitor

[Install]
WantedBy=multi-user.target